# leds.service
[Unit]
Description=Cool LEDs

[Service]
ExecStart=/usr/bin/python /home/pi/Documents/LEDMapper/oneColor.py

[Install]
WantedBy=multi-user.target

